# pi/yarbo-bridge-update.timer
[Unit]
Description=Check for Yarbo Bridge updates every 15 minutes

[Timer]
OnBootSec=2min
OnUnitActiveSec=15min
Persistent=true

[Install]
WantedBy=timers.target
